# lib/systemd/system/bluebinder.service.d/10-sleep.conf
[Service]
ExecStartPre=-/usr/sbin/rfkill unblock bluetooth
ExecStartPre=/bin/sleep 2
