# etc/systemd/system/android-mount.service.d/20-nologd.conf
[Service]
ExecStartPost=/usr/bin/mount --bind /dev/null /var/lib/lxc/android/rootfs/system/bin/logd
ExecStartPost=-(pgrep -f "/system/bin/logd" | xargs kill) || true
