# debian/adaptation-exynos9810-configs.getepoch.service
[Unit]
Description=set epoch on startup
After=lxc@android.service

[Service]
ExecStart=/usr/bin/getepoch

[Install]
WantedBy=multi-user.target
