# usr/lib/systemd/system/android-service@hwcomposer.service.d/40-exynos.conf
[Service]
ExecStartPre=/usr/bin/sleep 5
ExecStartPre=/usr/bin/mount -o bind /dev/null /vendor/lib64/hw/hwcomposer.exynos9810.so
